# systemd/netlink-notify.service
[Unit]
Description=Netlink Notification

[Service]
Type=simple
ExecStart=/usr/bin/netlink-notify

[Install]
WantedBy=default.target
